# weather-dashboard.service
[Unit]
Description=Weather Dashboard Web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/weather
ExecStart=/usr/bin/python3 /home/pi/weather/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
